# crontab
0 0 * * * supervisorctl restart mtproxy
